# deploy/supervisor/cbas.conf
; ============================================================================
; CBAS — Supervisor config
; Copy to:  /etc/supervisor/conf.d/cbas.conf
; Then:     sudo supervisorctl reread && sudo supervisorctl update
; Logs dir: sudo mkdir -p /var/log/cbas && sudo chown www-data:www-data /var/log/cbas
; ============================================================================

; ---- Next.js web app (served at digixworkspace.com/cbas via Nginx) ----------
[program:cbas-web]
directory=/var/www/cbas
command=/var/www/cbas/node_modules/.bin/next start -H 127.0.0.1 -p 3050
environment=NODE_ENV="production",NODE_OPTIONS="--no-deprecation"
user=www-data
autostart=true
autorestart=true
startsecs=10
startretries=3
stopsignal=SIGINT
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/cbas/web.out.log
stderr_logfile=/var/log/cbas/web.err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB

; ---- Python AI Engine (internal only, called by the web app) ----------------
[program:cbas-ai]
directory=/var/www/cbas/ai-engine
command=/var/www/cbas/ai-engine/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8088 --workers 1
environment=PYTHONUNBUFFERED="1"
user=www-data
autostart=true
autorestart=true
startsecs=10
startretries=3
stopasgroup=true
killasgroup=true
stdout_logfile=/var/log/cbas/ai.out.log
stderr_logfile=/var/log/cbas/ai.err.log
stdout_logfile_maxbytes=10MB
stderr_logfile_maxbytes=10MB
